# GraphMaker/crontab
0 7 * * * cd /code/graphmaker && python3 graph.py >> /var/log/cron.log 2>&1
